# crypto/oraclewebtoken/challenge/configs/supervisord.conf
[supervisord]
user=root
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:web]
user=www-data
group=www-data
directory=/app
command=/app/api-start.sh
environment=FLAG="ISS{oR4cL3s_pR3d1Ct3d_tH15_sT4nDaReD_w0vLd_f4Il}",SECRET_KEY="431e840a88a997a84873d074fd2871b5"
autostart=true
priority=1000
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0